# dragon_queen_server.service
[Unit]
Description=Zombrise 3D Game Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/zombrise_3d
ExecStart=/root/zombrise_3d/target/release/zombrise_server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits (optional)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
